# misc/configs/supervisor.conf
[program:apiserver]
directory=/path_to_project
command=/path_to_project/app
process_name=%(program_name)s
redirect_stderr=false
stdout_logfile=/path_to_project/logs/%(program_name)s.stdout.log
stderr_logfile=/path_to_project/logs/%(program_name)s.stderr.log
autorestart=true
stdout_logfile_maxbytes=100MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stdout_capture_maxbytes=100MB   ; number of bytes in 'capturemode' (default 0)
stdout_events_enabled=true   ; emit events on stdout writes (default false)
stderr_logfile_maxbytes=100MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
stderr_capture_maxbytes=100MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=true   ; emit events on stderr writes (default false)


[program:httplive]
directory=/path_to_httplive
command=/path_to_httplive/httplive -d /path_to_httplive/httplive.db
process_name=%(program_name)s
redirect_stderr=false
stdout_logfile=/path_to_httplive/logs/%(program_name)s.stdout.log
stderr_logfile=/path_to_httplive/logs/%(program_name)s.stderr.log
autorestart=true
